# src/etc/cron.hourly/update-packages
#!/bin/sh

{
	nice -n 19 pacman -Sqywu --noconfirm
} >/dev/null 2>&1

exit 0
